# dot_config/systemd/user/ical2org.timer
[Unit]
Description=Run ical2org timer

[Timer]
OnActiveSec=30s
OnUnitActiveSec=5m
Unit=ical2org.service

[Install]
WantedBy=timers.target
