# dot_config/systemd/user/ical2org.service
[Unit]
Description=Sync Owncloud + Todoist Calenders into OrgMode

[Service]
Type=oneshot
ExecStart=/home/bascht/Code/ical2org/sync.sh
