# dot_config/systemd/user/org-clock.timer
[Unit]
Description=Run org-clock timer
Wants=emacs.service

[Timer]
OnActiveSec=30s
OnUnitActiveSec=120m
Unit=org-clock.service

[Install]
WantedBy=timers.target