# dot_config/systemd/user/org-clock.service
[Unit]
Description=Write latest org-clock to the status file

[Service]
Type=oneshot
ExecStart=/home/bascht/bin/org-update-clock-file