# home/.config/systemd/user/maid.service
[Unit]
Description=Run maid cleanup

[Service]
Type=oneshot
ExecStart=maid clean --force