# home/.config/systemd/user/maid.timer
[Unit]
Description=Run maid timer

[Timer]
OnActiveSec=30s
OnUnitActiveSec=120m
Unit=maid.service

[Install]
WantedBy=timers.target