# home/.config/systemd/user/worklogger.service
[Unit]
Description=Ask for periodic worklog updates

[Service]
Type=oneshot
ExecStart=/home/bascht/bin/worklogger